# crontab
# custom crontab for dotcloud services
# runs the python news fetcher service in the background
0-59/15 * * * * python fetcher.py db.shlv
